# www/vaults3/files/vaults3.in
#!/bin/sh
#
# PROVIDE: vaults3
# REQUIRE: SERVERS DAEMON LOGIN
# BEFORE:    securelevel
# KEYWORD:   shutdown

. /etc/rc.subr

name="vaults3"
rcvar="vaults3_enable"

: ${vaults3_enable="NO"}
: ${vaults3_datadir="/var/db/vaults3"}
: ${vaults3_configfile="/usr/local/etc/${name}.yaml"}
: ${vaults3_user="vaults3"}
: ${vaults3_group="vaults3"}

command="/usr/local/bin/${name}"
pidfile="/var/run/${name}/${name}.pid"

# Assign explicit custom lifecycle actions
start_precmd="${name}_prestart"
start_cmd="${name}_start"

vaults3_prestart() {
    if [ ! -d "${vaults3_datadir}" ]; then
        mkdir -p "${vaults3_datadir}"
    fi
    chown -R ${vaults3_user}:${vaults3_group} "${vaults3_datadir}"

    if [ ! -d "/var/run/${name}" ]; then
        mkdir -p "/var/run/${name}"
    fi
    chown ${vaults3_user}:${vaults3_group} "/var/run/${name}"
}

vaults3_start() {
    echo "Starting ${name}..."

    touch "/var/log/${name}.log"
    chown "${vaults3_user}:${vaults3_group}" "/var/log/${name}.log"

    cd "${vaults3_datadir}" && \
    su -m "${vaults3_user}" -c "/usr/sbin/daemon -p ${pidfile} -o /var/log/${name}.log -f ${command} --config ${vaults3_configfile}"
}

run_rc_command "$1"
